# deploy/profile.d/verdi_complete.sh

# Verdi autocomplete
eval "$(_VERDI_COMPLETE=source verdi)"

